# docker/cron/crontab
* * * * * php yii payments/process-queue